# docker/provision/supervisord.conf
[supervisord]
nodaemon=true

[program:pmsiplan]
command=/start.sh
