# docker/service/pmsiplan/run
#!/bin/bash

cd /app/server
exec node /app/server/app/index.js
