# docker/service/pmsiplan/run
#!/bin/bash

cd /app/server && node ./app/index.js
